# touch-timeout.service
[Unit]
Description=Touchscreen Brightness & Timeout Manager
After=multi-user.target systemd-modules-load.service systemd-backlight@.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/touch-timeout
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/sys/class/backlight

# Resource management (removed Nice/IOSchedulingClass for touch responsiveness)

[Install]
WantedBy=multi-user.target
