# systemd/touch-timeout.service
[Unit]
Description=Touchscreen Brightness & Timeout Manager
Documentation=man:touch-timeout(1)
After=multi-user.target systemd-modules-load.service systemd-backlight@.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/touch-timeout
Restart=on-failure
RestartSec=5

# Logging configuration (reduce SD card writes)
StandardOutput=journal
StandardError=journal
SyslogLevel=warning

# Watchdog configuration (requires systemd support in binary)
# Daemon must send sd_notify("WATCHDOG=1") within this interval
# If watchdog timeout occurs, systemd will restart the service
# WatchdogSec=30s

# Timeout configuration
TimeoutStartSec=10s
# Allow 5 seconds for graceful shutdown
TimeoutStopSec=5s

# Security hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/sys/class/backlight

# Resource limits
# No Nice/IOScheduling restrictions for touch responsiveness
LimitNOFILE=64

[Install]
WantedBy=multi-user.target
