# systemd/touch-timeout.service
[Unit]
Description=Touchscreen Brightness & Timeout Manager
After=multi-user.target systemd-modules-load.service systemd-backlight@.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/touch-timeout
# Customize with systemd override: sudo systemctl edit touch-timeout
# Example override:
#   [Service]
#   ExecStart=
#   ExecStart=/usr/bin/touch-timeout -b 200 -t 600
Restart=on-failure
RestartSec=5

# Logging (quiet by default, add -v for verbose)
StandardOutput=journal
StandardError=journal
LogLevelMax=info

TimeoutStartSec=10s
TimeoutStopSec=5s

# Security hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/sys/class/backlight

LimitNOFILE=64

[Install]
WantedBy=multi-user.target
